# docker.service.d/40-docker.conf
[Unit]
Requires=flanneld.service
After=flanneld.service

[Service]
Environment="DOCKER_OPTS=--insecure-registry 192.168.30.149"
